# launchagents/am.dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>am.dispatcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/michaeloneal/.bun/bin/bun</string>
        <string>run</string>
        <string>/Users/michaeloneal/am/scripts/dispatcher.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/michaeloneal/am</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/michaeloneal/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>BOARD_URL</key>
        <string>http://localhost:4200</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/am-dispatcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/am-dispatcher.log</string>
</dict>
</plist>
